# build/cald.initd
#!/sbin/openrc-run

daemon="/usr/lib/cald/daemon.py"
command="/usr/lib/cald/calendar.py"
pidfile="/run/$RC_SVCNAME.pid"

start() {
    # Start daemon process
    python3 $daemon &
    # Save pid of daemon
    printf $! > $pidfile
}


stop() {
    # Stop the daemon if present
    kill -2 `cat $pidfile`
}
